# home/dot_config/systemd/user/swaync.service.d/20-enoshima.conf
[Service]
ExecStart=
ExecStart=%h/.local/bin/enoshima-swaync
