# home/dot_config/autostart/discord.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Discord (autostart)
Comment=Start Discord minimized with the graphical session
TryExec=/home/kentakang/.local/bin/discord-wayland
Exec=/home/kentakang/.local/bin/discord-wayland --start-minimized
Terminal=false
X-GNOME-Autostart-enabled=true
