# home/dot_config/systemd/user/cyberdock-event-bridge.service
[Unit]
Description=Synchronize native Hyprland minimize events with Cyberdock
PartOf=graphical-session.target
After=graphical-session.target cyberdock.service
ConditionEnvironment=WAYLAND_DISPLAY
ConditionEnvironment=HYPRLAND_INSTANCE_SIGNATURE

[Service]
Type=simple
ExecStart=%h/.local/bin/cyberdock-event-bridge
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target
