# home/dot_config/systemd/user/desktop-display-revert.service
[Unit]
Description=Revert an unconfirmed display mode transaction

[Service]
Type=oneshot
ExecStart=%h/.local/bin/desktop-display-mode revert
