# home/dot_config/systemd/user/desktop-display-events.service
[Unit]
Description=Reconcile saved display modes after Hyprland topology changes
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY
ConditionEnvironment=HYPRLAND_INSTANCE_SIGNATURE

[Service]
Type=simple
ExecStart=%h/.local/bin/desktop-display-event-listener
Restart=always
RestartSec=2
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target
