# home/dot_config/systemd/user/hyprsunset-quick.service
[Unit]
Description=Cyberpunk Library night-light toggle
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/hyprsunset --temperature 4500
Restart=on-failure
RestartSec=2s
TimeoutStopSec=5s
