# home/dot_config/systemd/user/kakaotalk-focus-guard.service
[Unit]
Description=Classify KakaoTalk XWayland surfaces and repair keyboard focus
PartOf=graphical-session.target
After=graphical-session.target cyberdock.service xembed-sni-proxy.service
ConditionEnvironment=WAYLAND_DISPLAY
ConditionEnvironment=HYPRLAND_INSTANCE_SIGNATURE

[Service]
Type=simple
ExecStart=%h/.local/bin/kakaotalk-focus-guard
Restart=always
RestartSec=2
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target
