# home/dot_config/systemd/user/protonmail-bridge.service
[Unit]
Description=Proton Mail Bridge local IMAP/SMTP service
Documentation=https://proton.me/support/bridge-for-linux
After=graphical-session.target
PartOf=graphical-session.target
ConditionPathExists=%S/protonmail-bridge/managed-service-enabled
StartLimitIntervalSec=5min
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/bin/protonmail-bridge-core --noninteractive
Restart=on-failure
RestartSec=20s
KillMode=process
UMask=0077
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=graphical-session.target
